# nut/rootfs/etc/services.d/upsmon/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Network UPS Tools
# Run upsmon
# ==============================================================================
if bashio::config.equals 'mode' 'netserver' ;then
    bashio::net.wait_for 3493
fi

bashio::log.info "Starting the UPS monitor and shutdown controller..."
exec upsmon -p -D
