# nut/rootfs/etc/services.d/upsd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Network UPS Tools
# Run upsd
# ==============================================================================
if bashio::config.equals 'mode' 'netserver' ;then
    bashio::log.info "Starting the UPS information server..."
    exec upsd -D
fi
